# etc/profile.d/cl_vrf.sh
VRF=$(/usr/cumulus/bin/cl-vrf identify prompt)
export VRF

PS1='\u@\h${VRF}:\w\$ '
